# roles/models/files/2004bay.conf
[program:{{ model }}]
command={{ venvabs }}/bin/mmi-runner {{ engine }} {{ configfile }} --pause
autostart=true
autorestart=true
startsecs=30
startretries=10
user={{ remote_user }}
directory=/var/models/{{ model }}

