# roles/models/files/sfobay.conf
[program:sfobay]
command=/home/ubuntu/venvs/mmi/bin/mmi-runner dflowfm r02b_bay.mdu --pause
autostart=true
autorestart=true
startsecs=30
startretries=10
user=ubuntu
directory=/var/models/sfobay

