# roles/mmi/templates/tracker.conf
[program:tracker]
command={{ venvabs }}/bin/mmi-tracker
autostart=true
autorestart=true
startsecs=30
startretries=10
user={{ ansible_ssh_user }}
directory=/home/{{ ansible_ssh_user }}

